# ops/track-articles.cron
# Update generated-articles.csv every 6 hours by scanning git log for auto-generated articles.
CRON_TZ=Asia/Jerusalem
0 */6 * * * cd /workspace/reginazaidler.github.io && /usr/bin/python3 scripts/track_generated_articles.py >> /workspace/reginazaidler.github.io/track-articles.log 2>&1
